# crontab
*/3 * * * * root /usr/sbin/update-ipsets >/dev/null 2>&1 && /sbin/netfilter-persistent save>/dev/null 2>&1
